# php/rootfs/etc/services.d/php-fpm/finish
#!/usr/bin/execlineb -S0

# bring down container if the services stop
s6-svscanctl -t /var/run/s6/services
