# steward/rootfs/etc/services.d/horizon/run
#!/usr/bin/with-contenv sh

exec -a laravel-horizon /usr/bin/php /var/www/artisan horizon >> /proc/self/fd/1 2>&1
