# beanstalkd/rootfs/etc/services.d/beanstalkd/run
#!/usr/bin/with-contenv sh

exec /usr/bin/beanstalkd -b /var/lib/beanstalkd/data >> /proc/self/fd/1 2>&1