# steward/rootfs/etc/services.d/scheduler/run
#!/usr/bin/execlineb -P

with-contenv

exec /usr/sbin/crond -f -c /etc/cron.d >> /proc/self/fd/1 2>&1

